# service-file/iotconnect-relay.service
[Unit]
Description=IoTConnect Relay Server
After=network.target

[Service]
Type=simple
User=<your-user>
WorkingDirectory=/opt/demo
ExecStart=/usr/bin/python3 /opt/demo/iotc-relay-server.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables (if needed)
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
